# sonos-avr.service
[Unit]
Description=Sonos AVR Controller
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sonos-avr
ExecStart=/opt/sonos-avr/sonos-avr.py
Restart=always
